# layout/Library/LaunchDaemons/ch.xxtou.hudservices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>DISABLE_TWEAKS</key>
		<string>1</string>
	</dict>
	<key>GroupName</key>
	<string>wheel</string>
	<key>HighPriorityIO</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>ch.xxtou.hudservices</string>
	<key>POSIXSpawnType</key>
	<string>App</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/var/jb/Applications/TrollSpeed.app/TrollSpeed</string>
		<string>-hud</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>UserName</key>
	<string>root</string>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>Managed</key>
			<false/>
			<key>Reported</key>
			<false/>
		</dict>
	</dict>
</dict>
</plist>
