# layout/Library/LaunchDaemons/ch.xxtou.hudapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>ch.xxtou.hudapp</string>
	<key>MachServices</key>
	<dict>
		<key>ch.xxtou.hudapp.springboard.launched</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/XXTAssistiveTouch.app/XXTAssistiveTouch</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>